# dircolors/dot-config/zsh/profile.d/80-dircolors.zsh
#(( $+commands[dircolors-gen] )) &&
#	[[ ! -f "$DIRCOLORS_INIT" ]] &&
#	dircolors-gen
